# services/serial-getty@ttyGS0.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- \u' -L 115200 %I vt100
