# etc/systemd/system/p2p-events.path
[Unit]
Description=Start p2p-events when wpa_supplicant ctrl iface for wlan0 appears

[Path]
PathExists=/var/run/wpa_supplicant/wlan0
Unit=p2p-events.service

[Install]
WantedBy=multi-user.target
