# etc/systemd/system/p2p-autoconnect.service
[Unit]
Description=Wi-Fi Direct Autoconnect (listen for GO-NEG-REQUEST and start WPS-PBC)
After=network.target p2p-wpa.service p2p-events.service
Requires=p2p-wpa.service
ConditionPathExists=/var/run/wpa_supplicant/p2p-dev-wlan0

[Service]
Type=simple
ExecStart=/bin/bash /home/receiveit/receiveit-server/scripts/p2p-autoconnect.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
